# supervisor-queue.conf
; ════════════════════════════════════════════════════
;  Supervisor — Worker de queue Laravel
;  Placer dans : /etc/supervisor/conf.d/ecommerce-worker.conf
;  Puis : supervisorctl reread && supervisorctl update
; ════════════════════════════════════════════════════

[program:ecommerce-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --timeout=90 --max-jobs=500 --max-time=3600
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/html/storage/logs/worker.log
stdout_logfile_maxbytes=10MB
stopwaitsecs=120
